# systemd/dev-agent-server.container
[Unit]
Description=dev-agent server (Hono + Claude agent loop)
Wants=network-online.target
After=network-online.target
Requires=podman.socket
Wants=dev-agent-proxy.service
After=dev-agent-proxy.service podman.socket

[Container]
ContainerName=dev-agent-server
Image=localhost/dev-agent-server:latest

# Secrets + tunables. Keep this file mode 0600.
EnvironmentFile=%h/dev-agent-server/.env

# Fixed paths inside the container.
Environment=WORKSPACE_DIR=/data/workspaces
Environment=DATA_DIR=/data/db
# Host path for the workspaces volume. The server uses this to translate
# per-session worktree paths before passing them to podman, because podman
# (running on the host) can't resolve the in-container path. The default
# location for a rootless-podman named volume is
#   $HOME/.local/share/containers/storage/volumes/<name>/_data
# If you've customized podman's storage, override this.
Environment=WORKSPACE_DIR_HOST=%h/.local/share/containers/storage/volumes/dev-agent-workspaces/_data
Environment=DOCKER_HOST=unix:///var/run/docker.sock
# podman-remote uses CONTAINER_HOST, not DOCKER_HOST. Set both so the server
# works regardless of which CLI is selected via ENGINE_CLI.
Environment=CONTAINER_HOST=unix:///var/run/docker.sock
# Use the podman CLI inside the container. It accepts the same flags as docker
# for our use cases AND supports --userns=keep-id (which the docker CLI rejects
# client-side). Required for rootless podman with SANDBOX_USERNS=keep-id.
Environment=ENGINE_CLI=podman
Environment=PROXY_URL=http://dev-agent-proxy:8888
# Network name for per-session sandbox containers. Matches NetworkName= in
# agent-egress.network. Quadlet only applies the systemd- prefix when
# NetworkName is NOT set; since we set it explicitly, the name is just
# agent-egress.
Environment=SANDBOX_NETWORK=agent-egress

Volume=dev-agent-workspaces.volume:/data/workspaces:Z
Volume=dev-agent-db.volume:/data/db:Z
# Rootless podman API socket → docker-compatible socket inside the container.
# %t = $XDG_RUNTIME_DIR (e.g. /run/user/1000).
# SecurityLabelDisable lets the container process connect to the socket without
# SELinux blocking it — necessary on Fedora where the container label would
# otherwise be denied access to the relabelled socket file.
SecurityLabelDisable=true
Volume=%t/podman/podman.sock:/var/run/docker.sock:z
# The seccomp profile is read by podman on the HOST (not inside the server
# container), so SECCOMP_PROFILE_HOST must be a host path. We also mount it
# into the server container so the existence-check in sandbox.ts can find it.
Volume=%h/dev-agent-server/sandbox/seccomp.json:/app/sandbox/seccomp.json:ro,Z
# The proxy's project.txt is bind-mounted here (rw) so the server can update
# it and then restart the proxy container to pick up the new allowlist.
Volume=%h/dev-agent-server/proxy/project.txt:/app/proxy/project.txt:rw,Z
Environment=PROXY_PROJECT_FILE=/app/proxy/project.txt
Environment=SECCOMP_PROFILE=/app/sandbox/seccomp.json
Environment=SECCOMP_PROFILE_HOST=%h/dev-agent-server/sandbox/seccomp.json

# The server only needs the default external network.
# It reaches GitHub and Anthropic directly, and reaches the proxy by container
# name (both are on the podman default network).
# It does NOT need to be on agent-egress -- it places sandbox containers onto
# that network via the podman API socket; it doesn't live there itself.
# (Being on agent-egress caused its DNS to break because that network is
# internal: true, which strips the upstream resolver.)
Network=podman

# Cloudflared (running as system service) connects to localhost:3737.
# Container internal port stays 3000 (the server's default); only the host-side
# published port is 3737. To change, edit just this line and the cloudflared
# tunnel config -- nothing inside the container needs to change.
PublishPort=127.0.0.1:3737:3000

[Service]
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
