# systemd/dev-agent-proxy.container
[Unit]
Description=dev-agent egress proxy (tinyproxy with allowlist)
Wants=network-online.target
After=network-online.target
Requires=agent-egress-network.service
After=agent-egress-network.service

[Container]
ContainerName=dev-agent-proxy
Image=localhost/dev-agent-proxy:latest
# Order matters: the first network's interface gets the default route inside
# the container. podman (external) must come first so outbound traffic to
# allowlisted hosts can actually leave the box. agent-egress (Internal=true)
# has no external route; if it were first, DNS and HTTPS to the internet
# would silently black-hole.
Network=podman
# Static IP so sandbox containers (which have no DNS on agent-egress) can
# reach the proxy by address. Must be within the Subnet= defined in
# agent-egress.network.
Network=agent-egress.network:ip=10.89.100.10
# Explicit DNS so the proxy can resolve external hostnames. When on both an
# internal network (agent-egress) and the default network, podman may pick the
# internal network's aardvark-dns as primary — which can't forward external
# queries. The proxy never needs to resolve container names, only external
# hosts, so a public resolver is correct here.
DNS=1.1.1.1
DNS=8.8.8.8
# %h is $HOME. Adjust if you place the repo elsewhere.
Volume=%h/dev-agent-server/proxy/project.txt:/etc/proxy/project.txt:ro,Z

[Service]
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
